# bin/counter.service
[Unit]
Description=counter service

[Service]
User=counter
WorkingDirectory=/usr/local/bin
ExecStart=/usr/bin/python3 /usr/local/bin/counter.py
# Optional items below
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target